# debian/midolman-agent-bin.upstart
description "Midolman Net Agent TCP server"
author "Ryu Ishimoto <ryu@midokura.jp>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

chdir /var/run

pre-start script
        mkdir -p /var/run/midolman
        chown root:root /var/run/midolman

        mkdir -p /var/lock/midolman
        chown root:root /var/lock/midolman
end script

exec su -c "midolman-net-agent --conf=/etc/midolman-agent.conf" root

